# dot_config/systemd/user/kanshi.service
[Unit]
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
ExecStart=/usr/bin/kanshi
Restart=on-failure
