# dot_config/systemd/user/swaybg.service
[Unit]
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
ExecStart=/usr/bin/swaybg -m fill -i "%h/Pictures/tomasz-filipek-59OA9ycJhpI-unsplash.jpg"
Restart=on-failure
