# dot_config/systemd/user/theme-switch.timer
[Unit]
Description="Trigger light/dark theme switching"

[Timer]
OnStartupSec=10s
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
