# dot_config/systemd/user/theme-switch.service
[Unit]
Description="Switch to light or dark theme"

[Service]
Type=oneshot
ExecStart=%h/.local/bin/theme-switch
