# dot_config/systemd/user/xdg-desktop-portal.service.d/override.conf
[Unit]
# Niri sets the DISPLAY env var dynamically
# and it seems that xdg-desktop-portal really hates not
# having it (it will not open links).
After=niri.service
