# dot_config/systemd/user/plasma-polkit-agent.service.d/override.conf
[Unit]
After=xwayland-satellite.service
